# docker/docker/mail/supervisord.conf
[supervisord]
nodaemon=true

[program:postfix]
command=/usr/sbin/postfix start-fg
priority=10

[program:dovecot]
command=/usr/sbin/dovecot -F
priority=20
